# systemd/wa-cli-backup.timer
[Unit]
Description=Daily backup of the WhatsApp CLI archive

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
